# polymarket-daily-refresh.service
[Unit]
Description=Polymarket Daily Wallet Refresh
After=network.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/opt/polymarket-bot
Environment="PATH=/opt/polymarket-bot/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/opt/polymarket-bot/.env
ExecStart=/opt/polymarket-bot/venv/bin/python3 /opt/polymarket-bot/daily_wallet_refresh.py
StandardOutput=journal
StandardError=journal

